# jdk8-ffmpeg/systemd/java.service
[Unit]
Description=springboot
After=network.target

[Service]
Type=forking
ExecStart=/opt/java/java_service.sh start
ExecRestart=/opt/java/java_service.sh start
ExecStop=/opt/java/java_service.sh stop


[Install]
WantedBy=multi-user.target
